# ansible/roles/proxy/templates/proxy.service.j2
[Unit]
Description=Coin-Ops Proxy Service (Go)
After=network.target redis-server.service
Wants=network.target

[Service]
Type=simple
User=cognitor-proxy
EnvironmentFile=/etc/cognitor/proxy.env
ExecStart=/opt/cognitor/proxy/proxy
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
