# ansible/templates/coin-history.service
[Unit]
Description=Coin-Ops History Service
After=network.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory=/home/{{ ansible_user }}/history-service
ExecStart=/home/{{ ansible_user }}/history-service/history-service
Environment=PORT=8001
Environment=CORS_ORIGINS=http://192.168.56.10
EnvironmentFile=/etc/coin-ops/history.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
